# rhasspy-speech/rootfs/etc/s6-overlay/s6-rc.d/rhasspy-speech/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start rhasspy-speech service
# ==============================================================================
cd /usr/src

flags=()

# VAD
if bashio::config.false 'vad'; then
    flags+=('--no-vad')
fi

# Speex
if bashio::config.true 'speex'; then
    flags+=('--speex')
fi

# Transcribers
if bashio::config.false 'grammar'; then
    flags+=('--no-grammar')
fi

if bashio::config.false 'grammar_streaming'; then
    flags+=('--no-grammar-streaming')
fi

if bashio::config.false 'arpa'; then
    flags+=('--no-arpa')
fi

if bashio::config.false 'arpa_streaming'; then
    flags+=('--no-arpa-streaming')
fi

if bashio::config.true 'arpa_rescore'; then
    flags+=('--arpa-rescore')
fi

# Misc
if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

exec .venv/bin/python3 -m wyoming_rhasspy_speech \
    --uri 'tcp://0.0.0.0:10300' \
    --tools-dir /usr/src/tools \
    --train-dir /share/rhasspy-speech/train \
    --models-dir /share/rhasspy-speech/models \
    --volume-multiplier "$(bashio::config 'volume_multiplier')" \
    --vad-threshold "$(bashio::config 'vad_threshold')" \
    --before-speech-seconds "$(bashio::config 'before_speech_seconds')" \
    --speex-noise-suppression "$(bashio::config 'speex_noise_suppression')" \
    --speex-auto-gain "$(bashio::config 'speex_auto_gain')" \
    --word-norm-distance-threshold "$(bashio::config 'word_norm_distance_threshold')" \
    --char-norm-distance-threshold "$(bashio::config 'char_norm_distance_threshold')" \
    --grammar-max-active "$(bashio::config 'grammar_max_active')" \
    --grammar-lattice-beam "$(bashio::config 'grammar_lattice_beam')" \
    --grammar-acoustic-scale "$(bashio::config 'grammar_acoustic_scale')" \
    --grammar-beam "$(bashio::config 'grammar_beam')" \
    --arpa-max-active "$(bashio::config 'arpa_max_active')" \
    --arpa-lattice-beam "$(bashio::config 'arpa_lattice_beam')" \
    --arpa-acoustic-scale "$(bashio::config 'arpa_acoustic_scale')" \
    --arpa-beam "$(bashio::config 'arpa_beam')" \
    --arpa-rescore-order "$(bashio::config 'arpa_rescore_order')" \
    --arpa-rescore-acoustic-scale "$(bashio::config 'arpa_rescore_acoustic_scale')" \
    --hass-token "${SUPERVISOR_TOKEN}" \
    --hass-websocket-uri 'ws://supervisor/core/websocket' \
    --hass-ingress \
    --web-server-host '0.0.0.0'  ${flags[@]}
