# rhasspy-speech/rootfs/etc/s6-overlay/s6-rc.d/rhasspy-speech/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start rhasspy-speech service
# ==============================================================================
cd /usr/src

flags=()

# VAD
if bashio::config.false 'vad'; then
    flags+=('--no-vad')
fi

# Speex
if bashio::config.true 'speex'; then
    flags+=('--speex')
fi

# Transcribers
if bashio::config.true 'arpa_rescore'; then
    flags+=('--arpa-rescore')
fi

# Misc
if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

exec .venv/bin/python3 -m wyoming_rhasspy_speech \
    --uri 'tcp://0.0.0.0:10300' \
    --tools-dir /usr/src/tools \
    --train-dir /share/rhasspy-speech/train \
    --models-dir /share/rhasspy-speech/models \
    --volume-multiplier "$(bashio::config 'volume_multiplier')" \
    --vad-threshold "$(bashio::config 'vad_threshold')" \
    --before-speech-seconds "$(bashio::config 'before_speech_seconds')" \
    --speex-noise-suppression "$(bashio::config 'speex_noise_suppression')" \
    --speex-auto-gain "$(bashio::config 'speex_auto_gain')" \
    --norm-distance-threshold "$(bashio::config 'norm_distance_threshold')" \
    --max-active "$(bashio::config 'max_active')" \
    --lattice-beam "$(bashio::config 'lattice_beam')" \
    --acoustic-scale "$(bashio::config 'acoustic_scale')" \
    --beam "$(bashio::config 'beam')" \
    --nbest "$(bashio::config 'nbest')" \
    --decode-mode "$(bashio::config 'decode_mode')" \
    --arpa-rescore-order "$(bashio::config 'arpa_rescore_order')" \
    --hass-token "${SUPERVISOR_TOKEN}" \
    --hass-websocket-uri 'ws://supervisor/core/websocket' \
    --hass-ingress \
    --web-server-host '0.0.0.0'  ${flags[@]}
