# whisper.cpp/rootfs/etc/s6-overlay/s6-rc.d/whisper.cpp/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start whisper.cpp service
# ==============================================================================
flags=()

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

exec python3 -m wyoming_whisper_cpp \
    --whisper-cpp-dir '/usr/share/wyoming-whisper-cpp/whisper.cpp' \
    --uri 'tcp://0.0.0.0:10300' \
    --data-dir /data \
    --download-dir /data \
    --model "$(bashio::config 'model')" \
    --language "$(bashio::config 'language')" \
    --beam-size "$(bashio::config 'beam_size')" \
    --audio-context-base "$(bashio::config 'audio_context_base')" ${flags[@]}
