# assist_microphone/rootfs/etc/s6-overlay/s6-rc.d/assist_microphone/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start service
# ==============================================================================
export PYTHONPATH=/usr/src

extra_args=()
if bashio::config.true 'debug_logging'; then
    extra_args+=('--debug')
fi

if bashio::config.true 'awake_sound'; then
    extra_args+=('--awake-sound' '/usr/src/sounds/awake.wav')
fi

if bashio::config.true 'done_sound'; then
    extra_args+=('--done-sound' '/usr/src/sounds/done.wav')
fi

if bashio::config.true 'udp_mic'; then
    extra_args+=('--udp-mic' 5000)
fi

if bashio::config.true 'udp_snd'; then
    extra_args+=('--udp-snd' 6055)
fi

exec python3 -m hass_satellite \
    'homeassistant' \
    "$(bashio::config 'token')" \
    --mic-device "$(bashio::config 'mic_device')" \
    --snd-device "$(bashio::config 'snd_device')" \
    --volume "$(bashio::config 'volume')" \
    --wake-buffer-seconds "$(bashio::config 'wake_buffer_seconds')" \
    --vad "$(bashio::config 'vad')" \
    --vad-model /usr/src/models/silero_vad.onnx ${extra_args[@]}
