# openwakeword/rootfs/etc/s6-overlay/s6-rc.d/openwakeword/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start openWakeWord service
# ==============================================================================
flags=()

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

if bashio::config.true 'save_audio'; then
    flags+=('--output-dir' '/share/openwakeword')
fi

exec python3 -m wyoming_openwakeword \
    --uri 'tcp://0.0.0.0:10400' \
    --model "$(bashio::config 'model')" \
    --threshold "$(bashio::config 'threshold')" \
    --trigger-level "$(bashio::config 'trigger_level')" \
    --noise-suppression "$(bashio::config 'noise_suppression')" \
    --auto-gain "$(bashio::config 'auto_gain')" ${flags[@]}
