# vosk/rootfs/etc/s6-overlay/s6-rc.d/vosk/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start vosk service
# ==============================================================================
flags=()

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

if bashio::config.true 'limit_sentences'; then
    flags+=('--limit-sentences')
fi

if bashio::config.true 'allow_unknown'; then
    flags+=('--allow-unknown')
fi

exec python3 -m wyoming_vosk \
    --uri 'tcp://0.0.0.0:10300' \
    --data-dir /data \
    --download-dir /data \
    --data-dir /share/vosk/models \
    --sentences-dir /share/vosk/sentences \
    --correct-sentences "$(bashio::config 'correct_sentences')" \
    --language "$(bashio::config 'preload_language')" \
    --preload-language "$(bashio::config 'preload_language')" ${flags[@]}
